# project/microservices/src/frontend/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/src/server
autorestart=true

[program:node_exporter]
command=/node_exporter
autorestart=truef
